# deploy/systemd/ceymail-dashboard.service
[Unit]
Description=CeyMail Dashboard (Next.js)
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=ceymail-mc
Group=ceymail-mc
WorkingDirectory=/opt/mission-control/apps/dashboard/.next/standalone
EnvironmentFile=-/var/lib/ceymail-mc/.env.local
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1

# Security hardening
# Note: ProtectSystem is NOT used because the dashboard's sudo-based operations
# need write access to /etc/postfix, /etc/dovecot, /etc/opendkim, /etc/roundcube,
# /etc/nginx, /etc/letsencrypt, /var/mail, /usr (apt), etc. The ReadWritePaths
# list would be so broad it negates the protection. Security boundary is the
# sudoers allowlist + nologin shell, not filesystem isolation.
# NoNewPrivileges is also not used — required for sudo during install.
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Data directories (config + backups)
ReadWritePaths=/opt/mission-control/apps/dashboard/.next/standalone/data /var/lib/ceymail-mc /var/backups/ceymail

[Install]
WantedBy=multi-user.target
